# postgres/rootfs/etc/s6-overlay/s6-rc.d/init-postgres/run
#!/command/execlineb -P
foreground { /usr/local/bin/init-postgres.sh }
