# n8n/rootfs/etc/s6-overlay/s6-rc.d/n8n/run
#!/command/execlineb
exec chpst -u n8n /usr/bin/start-n8n.sh
