# redis/rootfs/etc/s6-overlay/s6-rc.d/init-redis/run
#!/command/execlineb -P
foreground { /usr/bin/init-redis.sh }
