# open-webui/rootfs/etc/s6-overlay/s6-rc.d/openwebui/run
#!/command/execlineb
exec chpst -u openwebui /usr/bin/start-openwebui.sh
