# redis/rootfs/etc/s6-overlay/s6-rc.d/redis/run
#!/command/execlineb
exec chpst -u redis /usr/bin/redis-server /etc/redis/redis.conf
