# open-webui/rootfs/etc/s6-overlay/s6-rc.d/init-openwebui/run
#!/command/execlineb -P
foreground { /usr/bin/init-openwebui.sh }
