# n8n/rootfs/etc/s6-overlay/s6-rc.d/init-n8n/run
#!/command/execlineb -P
foreground { /usr/bin/init-n8n.sh }
