# postgres/rootfs/etc/s6-overlay/s6-rc.d/postgres/finish
#!/bin/bash
exec chpst -u postgres /usr/lib/postgresql/16/bin/pg_ctl -D /var/lib/postgresql/16/main -m fast stop
